# deploy/quadlet/vf-cmdb-frontend.container
[Unit]
Description=Virtualfactor CMDB React/Nginx frontend
Wants=vf-cmdb-backend.service
After=vf-cmdb-backend.service

[Container]
ContainerName=vf_cmdb_frontend
# Build this image first:  podman build -t localhost/vf_cmdb_frontend:latest ./frontend
Image=localhost/vf_cmdb_frontend:latest
Network=vf-cmdb.network
PublishPort=8080:8080

[Service]
Restart=always
TimeoutStartSec=120

[Install]
WantedBy=default.target
